# .config/systemd/user/mbsync.service
[Unit]
Description=Mailbox sync (mbsync)

[Service]
Type=oneshot
ExecStart=/usr/bin/mbsync -a
