# systemd-user/.config/systemd/user/battery-warn.service
[Unit]
Description=Battery warning

[Service]
Type=oneshot
Environment="DISPLAY=:0"
ExecStart=/home/brian/repos/scripts/battery-warn
