# systemd-user/.config/systemd/user/battery-warn.timer
[Unit]
Description=Battery warning timer

[Timer]
OnCalendar=minutely
Persistent=true
Unit=battery-warn.service

[Install]
WantedBy=timers.target
